# systemd/inika-nginx.service
[Unit]
Description=Inika Bot Nginx Reverse Proxy
After=network.target inika-backend.service inika-frontend.service

[Service]
Type=forking
ExecStart=/usr/sbin/nginx -c /opt/inika-bot/nginx/nginx.conf
ExecReload=/usr/sbin/nginx -s reload
ExecStop=/usr/sbin/nginx -s quit
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
