# systemd/inika-frontend.service
[Unit]
Description=Inika Bot Frontend (Next.js)
After=network.target inika-backend.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/inika-bot/frontend
Environment="NODE_ENV=production"
Environment="PORT=3000"
ExecStart=/usr/bin/node /opt/inika-bot/frontend/.next/standalone/server.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
